# scripts/launchagents/com.matias.ideas-graphiti-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matias.ideas-graphiti-ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mmirvois/clawd/venv/bin/python</string>
        <string>/Users/mmirvois/projects/ideas-pipeline/scripts/ingest_obsidian_to_graphiti.py</string>
        <string>--limit</string>
        <string>25</string>
    </array>

    <!-- Ingest promoted Obsidian notes → Graphiti KG.
         Nightly at 03:00 local. --limit 25 caps per-run cost (LLM calls).
         The state file deduplicates unchanged notes on subsequent runs. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mmirvois/clawd/logs/ideas-graphiti-ingest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mmirvois/clawd/logs/ideas-graphiti-ingest.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
